# examples/netbbs.service
# systemd unit for running a NetBBS node under Linux (design doc §2.1
# Tier 2, issue #82). Install as /etc/systemd/system/netbbs.service,
# adjust User/Group/WorkingDirectory/ExecStart for your install layout,
# then:
#
#   systemctl daemon-reload
#   systemctl enable --now netbbs
#
# This example assumes a dedicated, unprivileged "netbbs" system user
# and a virtualenv at /var/lib/netbbs/.venv (created with
# `python3 -m venv /var/lib/netbbs/.venv && /var/lib/netbbs/.venv/bin/pip
# install netbbs[ssh,web]`, or a plain `pip install netbbs[ssh,web]`
# into whatever interpreter you point ExecStart at) -- adjust the path
# if you installed differently. NetBBS never daemonizes itself (design
# doc §13.8): it runs in the foreground and expects the service
# supervisor to manage backgrounding/restart, which is exactly what
# Type=simple here does.

[Unit]
Description=NetBBS node
Documentation=https://github.com/Thiesi/NetBBS
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=netbbs
Group=netbbs
WorkingDirectory=/var/lib/netbbs
ExecStart=/var/lib/netbbs/.venv/bin/python -m netbbs --config /etc/netbbs/netbbs.toml

# NetBBS treats SIGTERM as "please shut down properly" (design doc
# §13.8): it warns connected users, waits up to the configured
# shutdown.graceful_delay_seconds (60s by default), then disconnects
# and exits cleanly. TimeoutStopSec must stay comfortably above that
# default -- lower it only if you've also lowered graceful_delay_seconds
# in netbbs.toml.
KillSignal=SIGTERM
TimeoutStopSec=90

Restart=on-failure
RestartSec=5

# Hardening (optional, safe defaults for a network-facing service with
# no need to touch anything outside its own state directory). Loosen
# ReadWritePaths if your database/blob storage/identity directory lives
# somewhere other than WorkingDirectory.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/netbbs

[Install]
WantedBy=multi-user.target
